# image/stage3/02-systemd-units/files/etc/systemd/system/mozilla-iot-gateway.check-for-update.timer
[Unit]
Description=Run the Mozilla WebThings Gateway update checker daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
